# crontab
SHELL=/bin/bash

0  * * * * root /opt/take-snapshot.sh   > /proc/$(cat /var/run/crond.pid)/fd/1 2>&1
15 * * * * root /opt/expire-snapshot.sh > /proc/$(cat /var/run/crond.pid)/fd/1 2>&1
#
