# crontab
*/2  * * * * /opt/take-snapshot.sh  > /proc/1/fd/1 2>/proc/1/fd/2
*/2 * * * * /opt/epxire-snapshot.sh > /proc/1/fd/1 2>/proc/1/fd/2
